# dbus-i2c/service/log/run
#!/bin/sh
exec 2>&1
exec multilog t s99999 n8 /var/log/dbus-i2c
